# scripts/multi-ai-chatapp-dev.service
[Unit]
Description=AI Sister — Ted's dev instance (cli mode, sisters.ted-h.com)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Multi-Ai-Chatapp/server
EnvironmentFile=/home/ubuntu/Multi-Ai-Chatapp/server/.env.dev
ExecStart=/usr/bin/node /home/ubuntu/Multi-Ai-Chatapp/server/dist/index.js
Restart=on-failure
RestartSec=5
Environment=HOME=/home/ubuntu
# Keep PATH so child CLIs (claude, codex, gemini) are found — dev still
# uses the CLI subprocess path.
Environment=PATH=/home/ubuntu/.nvm/versions/node/current/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
